# services.conf
[program:postgresql]
command=/etc/init.d/postgresql start
user=postgres
autorestart=true
